# supervisor/templates/server.ini
[program:php-fpm]
command = /usr/sbin/php-fpm --nodaemonize --fpm-config /etc/php-fpm.d/www.conf -g "daemon off;"
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld  --pid-file=/var/run/mysqld/mysqld.pid --log-error=/var/log/mysqld.log --user=root
autostart=true
autorestart=true
user=root



[program:nginx]
command=/usr/sbin/nginx -g "daemon off; error_log /dev/stderr info;"
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0